# daemons/gemma-4-26B-A4B-it-UD-Q4_K_M.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.llama.server</string>

    <key>UserName</key>
    <string>bender</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/llama-server</string>
      <string>-m</string>
      <string>/Users/bender/models/gemma-4-26B-A4B-it-UD-Q4_K_M.gguf</string>
      <string>--host</string>
      <string>0.0.0.0</string>
      <string>--port</string>
      <string>8080</string>
      <string>--ctx-size</string>
      <string>32768</string>
      <string>--flash-attn</string>
      <string>on</string>
      <string>--cache-type-k</string>
      <string>q8_0</string>
      <string>--cache-type-v</string>
      <string>q8_0</string>
      <string>--batch-size</string>
      <string>512</string>
      <string>--ubatch-size</string>
      <string>512</string>
      <string>--threads</string>
      <string>8</string>
      <string>--n-gpu-layers</string>
      <string>-1</string>
      <string>--temp</string>
      <string>0.6</string>
      <string>--top-p</string>
      <string>0.9</string>
      <string>--repeat-penalty</string>
      <string>1.03</string>
      <string>--repeat-last-n</string>
      <string>128</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/bender/logs/llm.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/bender/logs/llm.err</string>
  </dict>
</plist>
