# daemons/omlx.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.llama.server</string>

    <key>UserName</key>
    <string>bender</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/omlx</string>
      <string>serve</string>
      <string>--model-dir</string>
      <string>/Users/bender/models</string>
      <string>--host</string>
      <string>0.0.0.0</string>
      <string>--port</string>
      <string>8080</string>
      <string>--memory-guard</string>
      <string>aggressive</string>
      <string>--memory-guard-gb</string>
      <string>24</string>
      <string>--hot-cache-max-size</string>
      <string>4GB</string>
      <string>--no-hf-cache</string>
      <string>--paged-ssd-cache-dir</string>
      <string>/Users/bender/cache</string>
      <string>--max-concurrent-requests</string>
      <string>1</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/bender/logs/llm.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/bender/logs/llm.err</string>
  </dict>
</plist>
